# files/etc/confd/templates/confd_services.conf.tmpl
#This file was autogenerated by confd from consul
#{{ $host_name := (getenv "HOSTNAME")}}
#{{ $host_conf := ( print "/hosts/" $host_name "/services" ) }}
#{{ $tree := (print $host_conf "/")}}
#{{ range lsdir $tree }}
#    {{ $service := . }}
#    {{ $disabled := getv (print $tree $service "/disabled") "false"}}
#    {{ if eq $disabled "false" }}
#        {{ $conf_name := getv (print $tree $service "/conf_name") "default"}}
#        {{ $image_name := getv (print "/services/" $service "/" $conf_name "/image") (getv (print "/services/" $service "/default/image") (print "docker.welespay.ru/" $service)) }}
#        {{ $tag := getv (print "/services/" $service "/" $conf_name "/tag") (getv (print "/services/" $service "/default/tag") "current") }}
#        {{ $image := print $image_name ":" $tag }}
#        {{ $run_opts := getv (print "/services/" $service "/" $conf_name "/run_opts") ( getv (print "/services/" $service "/default/run_opts") "--rm")}}
#        {{ $network := getv (print "/services/" $service "/" $conf_name "/network") ( getv (print "/services/" $service "/default/network") "host")}}
#        {{ $volume_opts := getv (print "/services/" $service "/" $conf_name "/volume_opts") ( getv (print "/services/" $service "/default/volume_opts") (print "-v /striped/data/" $service ":/data -v /striped/logs/" $service ":/var/log"))}}
#        {{ $args := getv (print "/services/" $service "/" $conf_name "/args") ( getv (print "/services/" $service "/default/args") "")}}
#        {{ $env_opts := getv (print "/services/" $service "/" $conf_name "/env_opts") ( getv (print "/services/" $service "/default/env_opts") "")}}
#        {{ $autostart := getv (print "/services/" $service "/" $conf_name "/autostart") ( getv (print "/services/" $service "/default/autostart") "true")}}
#        {{ $autorestart := getv (print "/services/" $service "/" $conf_name "/autorestart") ( getv (print "/services/" $service "/default/autorestart") "unexpected")}}
#        {{ $stopsignal := getv (print "/services/" $service "/" $conf_name "/stopsignal") ( getv (print "/services/" $service "/default/stopsignal") "TERM")}}
#        {{ $stopasgroup := getv (print "/services/" $service "/" $conf_name "/stopasgroup") ( getv (print "/services/" $service "/default/stopasgroup") "false")}}
[program:{{ print $service }}]
command="/usr/docker/bin/service.sh"
autostart={{print $autostart}}
autorestart={{print $autorestart}}
environment=
    IMAGE="{{print $image}}",
    NAME="{{print $service}}",
    RUN_OPTS="{{print $run_opts}}",
    NETWORK="{{print $network}}",
    VOLUME_OPTS="{{print $volume_opts}}",
    HOSTNAME_OPTS="--hostname {{print $host_name}}",
    ENV_OPTS="{{print $env_opts}}",
    ARGS="{{print $args}}"
stopsignal={{print $stopsignal}}
stopwaitsecs=300
stopasgroup={{print $stopasgroup}}
redirect_stderr=true
stdout_logfile=/var/log/{{print $service}}.out        ; stdout log path, NONE for none; default AUTO
     {{end}}
{{end}}
